# model/stable/Scheduled Task Creation from XML in Public Folder.json
{
  "schemaVersion": "1.0",
  "id": "b6cff00e-cd32-43de-9285-6c238399c52f",
  "name": "Scheduled Task Creation from XML in Public Folder",
  "description": "Detects the creation of scheduled tasks using an XML file in C:\\Users\\Public Folder, a tactic identified in an APT campaign for maintaining persistence or executing malicious commands",
  "riskLevel": "medium",
  "scope": "endpoint",
  "rules": [
    {
      "id": "50c0810f-3e70-4e1a-8bfe-761e9376cb97",
      "name": "Scheduled Task Creation from XML in Public Folder",
      "score": 1
    }
  ],
  "threshold": {
    "targetScore": 1,
    "withinSeconds": 600,
    "dedupSeconds": 3600
  },
  "requiredProducts": [
    [
      "sao",
      "pao",
      "sds",
      "pds",
      "xes"
    ]
  ],
  "tags": [
    "mitre.t1053"
  ],
  "playbooks": [
    {
      "id": "bd6c9ef1-965d-4204-b9f2-d39003708f18",
      "name": "Default Playbook",
      "provider": "WORKBENCH",
      "type": "DEFAULT"
    },
    {
      "id": "f87e38bc-b0a9-4997-8b6a-d0762a7bcace",
      "name": "Generic file/ip/url/email investigation",
      "provider": "AIR",
      "type": "AUTOMATED_RESPONSE"
    }
  ],
  "visibleLevel": 1,
  "globalModel": true,
  "silent": false,
  "author": "Ahmed Kamal & Bahaa Yamany",
  "team": "MENA-IR",
  "date": "2025/01/02",
  "lastUpdate": "2025-03-17T12:16:46.000Z"
}