# model/stable/Suspicious Scheduled Task Creation from a Misnamed AutoIT Binary.json
{
  "schemaVersion": "1.0",
  "id": "2fdbd61d-19d4-4b1a-af0a-77b91067b839",
  "name": "Suspicious Scheduled Task Creation from a Misnamed AutoIT Binary",
  "description": "Detects potential malicious scheduled task creation from a suspicious misnamed AutoIT script which could lead to system compromise if not expected.",
  "riskLevel": "high",
  "scope": "endpoint",
  "rules": [
    {
      "id": "faddd518-7046-4388-82d9-28e16f0d3a76",
      "name": "Suspicious Scheduled Task Creation from a Misnamed AutoIT Binary",
      "score": 1
    }
  ],
  "threshold": {
    "targetScore": 1,
    "withinSeconds": 3600,
    "dedupSeconds": 3600
  },
  "requiredProducts": [
    [
      "sao",
      "pao",
      "sds",
      "xes",
      "pds"
    ]
  ],
  "tags": [
    "mitre.t1053",
    "mitre.t1053.005",
    "mitre.t1036.008",
    "mitre.t1036"
  ],
  "playbooks": [
    {
      "id": "bd6c9ef1-965d-4204-b9f2-d39003708f18",
      "name": "Default Playbook",
      "provider": "WORKBENCH",
      "type": "DEFAULT"
    },
    {
      "id": "f87e38bc-b0a9-4997-8b6a-d0762a7bcace",
      "name": "Generic file/ip/url/email investigation",
      "provider": "AIR",
      "type": "AUTOMATED_RESPONSE"
    }
  ],
  "visibleLevel": 1,
  "globalModel": true,
  "silent": false,
  "author": "John Paul Lim",
  "team": "TAD Service Ops",
  "date": "2024/10/17",
  "lastUpdate": "2024-10-28T06:13:55.000Z"
}